# roles/scala/templates/etc/profile.d/scala.sh
export SCALA_HOME={{ scala_home }}
export PATH=$SCALA_HOME/bin:$PATH